# templates/nrpe.service.erb
[Unit]
Description=Nagios Remote Program Executor
Documentation=http://www.nagios.org/documentation
Conflicts=nrpe.socket
Requires=network.target

[Install]
WantedBy=multi-user.target

[Service]
Type=forking
User=<%= @nagios_user %>
Group=<%= @nagios_user %>
EnvironmentFile=/etc/sysconfig/nrpe
ExecStart=/usr/sbin/nrpe -c /etc/nagios/nrpe.cfg -d $NRPE_SSL_OPT
Nice=-15
